# docs/tentacle-web.service
[Unit]
Description=Tentacle Pro — Web Dashboard (port 3000)
After=network.target postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/backend-tentacle-pro

# 生产环境变量，读取 EnvironmentFile
EnvironmentFile=/root/backend-tentacle-pro/.env

ExecStart=/root/.bun/bin/bun run web/src/index.ts

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-web

[Install]
WantedBy=multi-user.target