# docs/tentacle-api.service
[Unit]
Description=Tentacle Pro — External API Service (port 3001)
After=network.target postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/backend-tentacle-pro/api

ExecStart=/root/.bun/bin/bun --env-file /root/backend-tentacle-pro/.env src/index.ts

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-api

[Install]
WantedBy=multi-user.target